# scripts/wdfwatch-token-refresh.timer
[Unit]
Description=WDFWatch OAuth Token Refresh Timer
Requires=wdfwatch-token-refresh.service

[Timer]
# Run every 30 minutes
OnCalendar=*:0/30
# Run 5 minutes after boot
OnBootSec=5min
# Randomize by up to 5 minutes to avoid exact timing
RandomizedDelaySec=5min
# If missed (e.g., system was off), run immediately
Persistent=true

[Install]
WantedBy=timers.target